# job-agent.service
[Unit]
Description=Dice Job Application Agent
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/job-agent
Environment="PATH=/home/ubuntu/job-agent/venv/bin"
ExecStart=/home/ubuntu/job-agent/venv/bin/python main.py start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
